# librechat/rootfs/etc/services.d/librechat/finish
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Service LibreChat finish
# ==============================================================================

bashio::log.info "Arrêt de LibreChat..."

# Graceful shutdown de Node.js
if pidof node > /dev/null; then
    pkill -TERM node
    sleep 5
fi

# Force kill si nécessaire
if pidof node > /dev/null; then
    pkill -KILL node
fi

bashio::log.info "LibreChat arrêté" 